# systemd/user/ags.service
[Unit]
Description=ags

[Service]
ExecStart=/usr/bin/ags run --gtk 4
Restart=always

[Install]
WantedBy=default.target
